# automation/launchd/com.hyejinha.team-dashboard.excel-autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hyejinha.team-dashboard.excel-autopush</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/hyejinha/.local/bin/team-dashboard-auto-push-excel.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/hyejinha</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/hyejinha</string>
        <key>REPO_ROOT</key>
        <string>/Users/hyejinha/Desktop/Workspace/Team</string>
        <key>RUN_ONCE</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>WatchPaths</key>
    <array>
        <string>/Users/hyejinha/Desktop/Workspace/Team/2026_멀티.xlsx</string>
    </array>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/hyejinha/Library/Logs/team-dashboard-excel-autopush.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/hyejinha/Library/Logs/team-dashboard-excel-autopush.err.log</string>
</dict>
</plist>
